# conf_files/new_system/start_ipacm_perf_le.service
# Automatically generated by systemd-sysv-generator

[Unit]
Documentation=man:systemd-sysv-generator(8)
SourcePath=/etc/init.d/start_ipacm_perf_le

After=populate-volatile.service
Before=devpts.service
Before=multi-user.target

[Service]
Type=forking
Restart=no
TimeoutSec=5min
IgnoreSIGPIPE=no
KillMode=process
GuessMainPID=no
RemainAfterExit=yes
ExecStart=/etc/init.d/start_ipacm_perf_le start
ExecStop=/etc/init.d/start_ipacm_perf_le stop
