# conf_files/new_system/start_qti_le.service
# Automatically generated by systemd-sysv-generator

[Unit]
Documentation=man:systemd-sysv-generator(8)
SourcePath=/etc/init.d/start_qti_le
After=start_QCMAP_ConnectionManager_le.service
Before=bootmisc.service
Before=multi-user.target

[Service]
Type=forking
Restart=no
TimeoutSec=5min
IgnoreSIGPIPE=no
KillMode=process
GuessMainPID=no
RemainAfterExit=yes
ExecStart=/etc/init.d/start_qti_le start
ExecStop=/etc/init.d/start_qti_le stop
